# debian/qemu-block-extra.postinst
#! /bin/sh
set -e

if [ "$1" = configure -a -x /usr/bin/deb-systemd-helper ] &&
   dpkg --compare-versions -- "$2" lt-nl 1:9.0.2+ds-1ubuntu1~
then
   deb-systemd-helper purge 'run-qemu.mount' >/dev/null || :
fi

#DEBHELPER#
